# sample.crontab
# Sample crontab(5). Submit sensor values every 5 minutes
*/5 * * * * lmsensors_cosm_submit.py -f cosm.cfg -l cosm.log
